# package/openwrt-captive-monitor/files/etc/init.d/captive-monitor
#!/bin/sh /etc/rc.common
# OpenWrt init script for Captive Portal Monitor
# shellcheck shell=sh
# shellcheck disable=SC2034 # START, STOP and USE_PROCD are consumed by rc.common metadata.

START=99
STOP=10

USE_PROCD=1

SCRIPT_PATH="/usr/sbin/openwrt_captive_monitor"

# shellcheck disable=SC1091 # Provided by OpenWrt base-files runtime
. /lib/functions.sh

load_config() {
    enabled=0
    mode="monitor"
    wifi_interface=""
    wifi_logical=""
    monitor_interval=""
    lan_interface=""
    lan_ip=""

    config_load captive-monitor
    config_get_bool enabled config enabled 0
    config_get mode config mode monitor
    config_get wifi_interface config wifi_interface ""
    config_get wifi_logical config wifi_logical ""
    config_get monitor_interval config monitor_interval ""
    config_get lan_interface config lan_interface ""
    config_get lan_ip config lan_ip ""
}

start_service() {
    if [ ! -x "$SCRIPT_PATH" ]; then
        logger -t captive-monitor -p user.err "Script not found or not executable: $SCRIPT_PATH"
        return 1
    fi

    load_config

    if [ "$enabled" -ne 1 ]; then
        logger -t captive-monitor -p user.info "Service disabled in UCI config, not starting"
        return 0
    fi

    mode_arg=""
    case "$mode" in
        oneshot)
            mode_arg="--oneshot"
            ;;
        *)
            mode="monitor"
            mode_arg="--monitor"
            ;;
    esac

    logger -t captive-monitor -p user.info "Starting Captive Portal Monitor (mode: $mode)"

    procd_open_instance
    procd_set_param command "$SCRIPT_PATH" "$mode_arg"

    if [ "$mode" = "monitor" ]; then
        procd_set_param respawn 3600 5 5
    fi

    procd_set_param stdout 1
    procd_set_param stderr 1

    [ -n "$wifi_interface" ] && procd_set_param env WIFI_INTERFACE="$wifi_interface"
    [ -n "$wifi_logical" ] && procd_set_param env WIFI_LOGICAL="$wifi_logical"
    [ -n "$monitor_interval" ] && procd_set_param env MONITOR_INTERVAL="$monitor_interval"
    [ -n "$lan_interface" ] && procd_set_param env LAN_INTERFACE="$lan_interface"
    [ -n "$lan_ip" ] && procd_set_param env LAN_IP="$lan_ip"

    procd_close_instance
}

stop_service() {
    logger -t captive-monitor -p user.info "Stopping Captive Portal Monitor"
}

reload_service() {
    stop
    start
}

service_triggers() {
    procd_add_reload_trigger "network" "captive-monitor"
}
