# init.d/captive-monitor
#!/bin/sh /etc/rc.common
# OpenWRT init script for Captive Portal Monitor
# /etc/init.d/captive-monitor

START=99
STOP=10

USE_PROCD=1

SCRIPT_PATH="/usr/bin/openwrt_captive_monitor.sh"
WIFI_INTERFACE="${WIFI_INTERFACE:-phy1-sta0}"
WIFI_LOGICAL="${WIFI_LOGICAL:-wwan}"
MONITOR_INTERVAL="${MONITOR_INTERVAL:-60}"

start_service() {
    # Проверка существования скрипта
    if [ ! -x "$SCRIPT_PATH" ]; then
        logger -t captive-monitor -p user.err "Script not found or not executable: $SCRIPT_PATH"
        return 1
    fi
    
    logger -t captive-monitor -p user.info "Starting Captive Portal Monitor"
    
    procd_open_instance
    procd_set_param command "$SCRIPT_PATH" --monitor \
        --interface "$WIFI_INTERFACE" \
        --logical "$WIFI_LOGICAL" \
        --interval "$MONITOR_INTERVAL"
    
    # Автоматический перезапуск при падении
    # threshold, timeout, retry
    procd_set_param respawn 3600 5 5
    
    # Логирование stdout/stderr
    procd_set_param stdout 1
    procd_set_param stderr 1
    
    # Переменные окружения
    procd_set_param env WIFI_INTERFACE="$WIFI_INTERFACE"
    procd_set_param env WIFI_LOGICAL="$WIFI_LOGICAL"
    procd_set_param env MONITOR_INTERVAL="$MONITOR_INTERVAL"
    
    procd_close_instance
}

stop_service() {
    logger -t captive-monitor -p user.info "Stopping Captive Portal Monitor"
    
    # Удаляем все правила captive portal при остановке сервиса
    
    # HTTP/HTTPS редирект
    iptables -t nat -D PREROUTING -i "$WIFI_INTERFACE" -j CAPTIVE_REDIRECT 2>/dev/null
    iptables -t nat -F CAPTIVE_REDIRECT 2>/dev/null
    iptables -t nat -X CAPTIVE_REDIRECT 2>/dev/null
    
    # DNS редирект
    iptables -t nat -D PREROUTING -i "$WIFI_INTERFACE" -j CAPTIVE_DNS_REDIRECT 2>/dev/null
    iptables -t nat -F CAPTIVE_DNS_REDIRECT 2>/dev/null
    iptables -t nat -X CAPTIVE_DNS_REDIRECT 2>/dev/null
    
    # DNS spoofing конфигурация
    rm -f /tmp/dnsmasq.d/captive-portal.conf
    /etc/init.d/dnsmasq restart 2>/dev/null || /etc/init.d/dnsmasq reload 2>/dev/null
    
    logger -t captive-monitor -p user.info "Captive Portal режим полностью удален"
}

reload_service() {
    stop
    start
}

service_triggers() {
    procd_add_reload_trigger "network"
}